# security/etc/uiam.service
# Service file for uIAM daemon
#
# Date: Oct 16, 2018
# Author: Andrea Di Cesare <andrea@softinstigate.com>

[Unit]
Description=uIAM Daemon
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/java -server -jar /usr/share/uiam/uiam.jar /etc/uiam/uiam.yml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
